# services/wifiFallback.service
[Unit]
Description=wireless-lab WiFi Fallback AP
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/piradio/wireless-lab
ExecStart=/home/piradio/wireless-lab/venv/bin/python /home/piradio/wireless-lab/services/wifiFallback.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
